# printserve-boot.service
[Unit]
Description=VESYL Print — boot splash
# Start early in boot and show a splash while the system comes up.
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target shutdown.target
Conflicts=shutdown.target printserve-display.service

[Service]
Type=simple
# Runs as root so it works before any user session exists; root can write
# /dev/fb1 and read the app files regardless of ownership.
ExecStart=/usr/bin/python3 /home/vesyl/printserve/main.py --boot
Restart=no

[Install]
WantedBy=sysinit.target
